# 60-clevo-airplane-mode.rules
KERNEL=="rfkill*", ACTION=="change", RUN+="/bin/bash /usr/bin/clevo-airplane-mode-led-control $number"
